# pps/crontab/pre_1day.sh
#!/bin/sh
source /home/httpd/ppysq.pt.pps.tv/crontab/function.sh